# roles/monitor/files/reboot-check.service
[Unit]
Description="System reboot check"
Requires=reboot-check.timer

[Service]
Type=simple
ExecStart=/usr/local/bin/reboot-check.sh
